# scripts/cctv-backend.service
# NOTE FOR USERS/DEVELOPERS:
# 'User=aryam' and '/home/aryam/Streamify/backend' are template placeholders.
# The 'install-services.sh' and 'setup-pi.sh' scripts will automatically detect
# your system username and repository path, replacing these values on installation.
# If manual setup is preferred, edit 'User' and 'WorkingDirectory' below.

[Unit]
Description=CCTV NVR System API Backend
After=network.target

[Service]
Type=simple
User=aryam
WorkingDirectory=/home/aryam/Streamify/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
